# stcm.service
[Unit]
Description=SillyTavern Campaign Manager
After=network.target ollama.service

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/stcm
ExecStart=/usr/bin/python3 backend/main.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
